# server/django.service
[Service]
Type=simple
ExecStart=/bin/bash /code/server/run.sh